# scripts/server/files/crontab
*/5 * * * * bash -lc "$HOME/scripts/timed_cleanup"
